# files/base/usr/lib/systemd/user/flatpak-update.timer
[Unit]
Description="Update Flatpaks"

[Timer]
OnUnitActiveSec=8h
Unit=flatpak-update.service

[Install]
WantedBy=multi-user.target