# files/base/usr/lib/systemd/user/flatpak-update.service
[Unit]
Description="Update Flatpaks"

[Service]
ExecStart=/usr/bin/flatpak update
ExecStart=/usr/bin/flatpak --user update
