# files/base/usr/lib/systemd/system/uupd.timer
[Unit]
Description="Run uupd"

[Timer]
OnBootSec=5min
OnUnitActiveSec=8h
Unit=uupd.service

[Install]
WantedBy=timers.target